# .history/docker/supervisord_20250308105721.conf
[supervisord]
nodaemon=true
user=documind
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:documind]
command=node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/app/logs/app-error.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/app/logs/app-output.log
stdout_logfile_maxbytes=10MB
environment=NODE_ENV=production

[program:cron-jobs]
command=node /app/dist/scheduler.js
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/cron-error.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/app/logs/cron-output.log
stdout_logfile_maxbytes=10MB
environment=NODE_ENV=production
